# .config/systemd/user/picom.service
[Unit]
Description=Picom compositor for X windows

[Service]
ExecStart=/usr/bin/picom
