# .config/systemd/user/picom.service
[Unit]
Description=Polybar status bar

[Service]
ExecStart=/usr/bin/picom
