# .config/systemd/user/polybar@.service
[Unit]
Description="Polybar status bar (%i)"

[Service]
ExecStart=/usr/bin/polybar %i
